# dot_config/private_systemd/private_user/mbsync.service
[Unit]
Description=Mailbox synchronization service

[Service]
Type=oneshot
ExecStart=/home/tbjers/bin/mail/syncmail
