# dot_config/private_systemd/private_user/mbsync.service.tmpl
[Unit]
Description=Mailbox synchronization service
RefuseManualStart=no
RefuseManualStop=yes

[Service]
Type=oneshot
ExecStart=/home/{{ .chezmoi.username }}/bin/mail/syncmail
