# dot_config/private_systemd/private_user/taskwarrior.service
[Unit]
Description=Taskwarrior synchronization service
RefuseManualStart=no
RefuseManualStop=yes

[Service]
Type=oneshot
ExecStart=task sync
