# dot_config/private_systemd/private_user/vdirsyncer.timer
[Unit]
Description=Synchronize vdirs

[Timer]
OnBootSec=5m
OnUnitActiveSec=5m
AccuracySec=1m

[Install]
WantedBy=timers.target
