# dot_config/autostart/ssh-add.desktop
[Desktop Entry]
Exec=ssh-add -q ~/.ssh/keys.private/*
Name=ssh-add
Type=Application
X-KDE-AutostartScript=true
